# deploy/systemd/workdir.service
[Unit]
Description=workdir — low-cost Firecracker microVM sandbox provider
Documentation=https://workdir.dev
After=network-online.target workdir-net.service
Wants=network-online.target workdir-net.service

[Service]
Type=simple
User=workdir
Group=workdir
# KVM access for Firecracker, CAP_NET_ADMIN to create/attach per-VM tap devices.
SupplementaryGroups=kvm
AmbientCapabilities=CAP_NET_ADMIN
CapabilityBoundingSet=CAP_NET_ADMIN
Environment=WORKDIR_CONFIG=/etc/workdir/config.toml
Environment=RUST_LOG=info,sandboxd=info
# --config is a global flag (before the subcommand); the unit passes it via the
# WORKDIR_CONFIG env var above, so `serve` takes no extra args.
ExecStart=/usr/local/bin/workdir serve
Restart=on-failure
RestartSec=2
LimitNOFILE=1048576
NoNewPrivileges=true
ReadWritePaths=/var/lib/workdir /etc/workdir
StateDirectory=workdir
# Firecracker + per-VM networking.
DeviceAllow=/dev/kvm rw
DeviceAllow=/dev/net/tun rw

[Install]
WantedBy=multi-user.target
