# deploy/systemd/workdir-backup.timer
[Unit]
Description=Run the workdir state backup nightly

[Timer]
OnCalendar=*-*-* 03:30:00
Persistent=true
RandomizedDelaySec=600

[Install]
WantedBy=timers.target
